# systemd/beam.service
[Unit]
Description=Beam Remote Desktop Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/beam-server --config /etc/beam/beam.toml
Restart=always
RestartSec=5
LimitNOFILE=65536

# Security hardening
# NOTE: ProtectSystem and ProtectHome must not be strict/yes because
# the server spawns agent processes as real users that need:
# - Home directory access (desktop config, file transfers)
# - /var/log write access (Xorg logs when using setuid wrapper)
# - /run/user access (PulseAudio, D-Bus)
ProtectSystem=full
ProtectHome=no
NoNewPrivileges=no
PrivateTmp=no

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=beam

[Install]
WantedBy=multi-user.target
